# nexus-webhook.service
[Unit]
Description=Nexus Webhook Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/git/ghabs/nexus-arc/examples/nexus-bot
EnvironmentFile=/home/ubuntu/git/ghabs/nexus/.env
Environment="PYTHONUNBUFFERED=1"
# Include virtualenv and system paths so `gh` and other system binaries are available
Environment="PATH=/home/ubuntu/git/ghabs/nexus/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/ubuntu/git/ghabs/nexus-arc"
ExecStart=/home/ubuntu/git/ghabs/nexus/venv/bin/python src/webhook_server.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/nexus/webhook.log
StandardError=append:/var/log/nexus/webhook_error.log

[Install]
WantedBy=multi-user.target
