# nexus-webhook.service
[Unit]
Description=Nexus GitHub Webhook Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/git/ghabs/nexus-core/examples/nexus-bot
ExecStart=/home/ubuntu/git/ghabs/nexus/venv/bin/python /home/ubuntu/git/ghabs/nexus-core/examples/nexus-bot/src/webhook_server.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/nexus/webhook.log
StandardError=append:/var/log/nexus/webhook_error.log

# Environment
Environment="PATH=/home/ubuntu/git/ghabs/nexus/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/ubuntu/git/ghabs/nexus/.env

[Install]
WantedBy=multi-user.target
