# nexus-webhook.service
[Unit]
Description=Nexus GitHub Webhook Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/git/ghabs/nexus
ExecStart=/home/ubuntu/git/ghabs/nexus/venv/bin/python3 /home/ubuntu/git/ghabs/nexus/src/webhook_server.py
Restart=always
RestartSec=10
StandardOutput=append:/home/ubuntu/git/ghabs/nexus/logs/webhook.log
StandardError=append:/home/ubuntu/git/ghabs/nexus/logs/webhook_error.log

# Environment
Environment="PATH=/home/ubuntu/git/ghabs/nexus/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/ubuntu/git/ghabs/nexus/vars.secret

[Install]
WantedBy=multi-user.target
